# supervisord_app.conf
[unix_http_server]
file=/opt/search_platform/tmp/supervisor.sock                       ; path to your socket file
chmod=0777                 ; socket file mode (default 0700)
sockchown=search_platform:search_platform       ; socket file uid:gid owner
# username=user              ; (default is no username (open server))
# password=123               ; (default is no password (open server))

[inet_http_server]
# Web管理界面设定
port=9001 ;
# username = user ;
# password = 1234 ;

[supervisord]
logfile=/opt/search_platform/logs/supervisor/supervisord.log                    ; supervisord log file
logfile_maxbytes=10MB                           ; maximum size of logfile before rotation
logfile_backups=1                              ; number of backed up logfiles
loglevel=info                                   ; info, debug, warn, trace
pidfile=/opt/search_platform/tmp/supervisord.pid                    ; pidfile location
nodaemon=false                                  ; run supervisord as a daemon
minfds=1024                                     ; number of startup file descriptors
minprocs=200                                    ; number of process descriptors
user=search_platform                                   ; default user
 
[rpcinterface:supervisor]
supervisor.rpcinterface_factory = supervisor.rpcinterface:make_main_rpcinterface

[supervisorctl]
serverurl=unix:///opt/search_platform/tmp/supervisor.sock           ; use a unix:// URL  for a unix socket
# username=user
# password=123

[program:Webserver_8201]
command= export LD_PRELOAD="/usr/lib64/libtcmalloc.so.4" &&  /home/search_platform/.virtualenvs/search_platform/bin/uwsgi --ini uwsgi/uwsgi_1.ini   ;
autostart=true  ;
stdout_logfile= /opt/search_platform/logs/supervisor/uwsgi_1.log  ;
process_name=%(program_name)s  ;
directory=/opt/search_platform ;
redirect_stderr=true ;
stopsignal=QUIT ;
reload-on-rss=300;
max-requests = 10000;
priority=1 ;

[program:Webserver_8202]
command= export LD_PRELOAD="/usr/lib64/libtcmalloc.so.4" &&  /home/search_platform/.virtualenvs/search_platform/bin/uwsgi --ini uwsgi/uwsgi_2.ini   ;
autostart=true  ;
stdout_logfile= /opt/search_platform/logs/supervisor/uwsgi_2.log  ;
process_name=%(program_name)s  ;
directory=/opt/search_platform ;
redirect_stderr=true ;
stopsignal=QUIT ;
reload-on-rss=300;
max-requests = 10000;
priority=1 ;

[program:Webserver_8203]
command= export LD_PRELOAD="/usr/lib64/libtcmalloc.so.4" &&  /home/search_platform/.virtualenvs/search_platform/bin/uwsgi --ini uwsgi/uwsgi_3.ini   ;
autostart=true  ;
stdout_logfile= /opt/search_platform/logs/supervisor/uwsgi_3.log  ;
process_name=%(program_name)s  ;
directory=/opt/search_platform ;
redirect_stderr=true ;
stopsignal=QUIT ;
reload-on-rss=300;
max-requests = 10000;
priority=1 ;

[program:Webserver_8204]
command= export LD_PRELOAD="/usr/lib64/libtcmalloc.so.4" &&  /home/search_platform/.virtualenvs/search_platform/bin/uwsgi --ini uwsgi/uwsgi_4.ini   ;
autostart=true  ;
stdout_logfile= /opt/search_platform/logs/supervisor/uwsgi_4.log  ;
process_name=%(program_name)s  ;
directory=/opt/search_platform ;
redirect_stderr=true ;
stopsignal=QUIT ;
reload-on-rss=300;
max-requests = 10000;
priority=1 ;

[program:Nginx]
command=/usr/sbin/nginx    ;
autostart=true  ;
stdout_logfile= /opt/search_platform/logs/supervisor/nginx.log  ;
process_name=%(program_name)s  ;
redirect_stderr=true ;
stopsignal=QUIT ;